# host/systemd/orion-netsec-nic-tuning.service
[Unit]
Description=Orion Sentinel NetSec - NIC Offload Tuning for Packet Capture
Documentation=https://github.com/orionsentinel/Orion-sentinel-netsec-ai
After=network-online.target orion-netsec-promisc.service
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes

# Disable NIC offload features for better packet capture reliability
# This can help reduce packet drops on high-traffic mirror ports
# Trade-off: Slightly higher CPU usage on the Pi 5
#
# Only enable this if you experience packet drops in ntopng or Suricata
# Check drops with: sudo ethtool -S eth0 | grep drop

# GRO (Generic Receive Offload) - can coalesce packets before capture
ExecStart=/usr/sbin/ethtool -K eth0 gro off

# LRO (Large Receive Offload) - similar to GRO
ExecStart=/usr/sbin/ethtool -K eth0 lro off

# GSO (Generic Segmentation Offload) - can affect transmitted packet visibility
ExecStart=/usr/sbin/ethtool -K eth0 gso off

# TSO (TCP Segmentation Offload) - similar to GSO
ExecStart=/usr/sbin/ethtool -K eth0 tso off

# Note: These settings are reset on interface down/up
# This service ensures they persist across reboots

# Restart policy
Restart=no

[Install]
WantedBy=multi-user.target
